# deploy/systemd/tun-runtime-server-watchdog@.service
[Unit]
Description=Tun Runtime Server Watchdog (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/tun-runtime-link-watchdog.sh --role server --instance %i

