# deploy/systemd/tun-runtime-client-watchdog@.service
[Unit]
Description=Tun Runtime Client Watchdog (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/tun-runtime-link-watchdog.sh --role client --instance %i

